# deploy/monitoring/cron/ai-tutor-index-cron.cron
# Sprint 4.1.4 — RAG indexing cron.
# Запускать каждый день в 03:30 UTC (после backup в 03:00, до maintenance).
30 3 * * * root /opt/ai-tutor/deploy/monitoring/index-cron.sh
