# deploy/monitoring/cron/ai-tutor-weekly-summary.cron
# Sprint 9.1 — еженедельный weekly summary для родителей (воскресенье 18:00 MSK = 15:00 UTC)
# Sprint 6.3 fix: source /opt/ai-tutor/.env (НЕ /etc/ai-tutor/.env — он не существует на проде)
0 15 * * 0 root bash -c "set -a; source /opt/ai-tutor/.env; set +a; docker exec -u root deploy-backend-1 python3 -m app.scripts.weekly_summary" >> /var/log/ai-tutor-weekly-summary.log 2>&1
