# deploy/monitoring/cron/ai-tutor-backup-verify.cron
# Sprint 3.6.1 — еженедельная верификация backup через test-restore.
# Понедельник 04:00 UTC.
0 4 * * 1 root bash -c "bash /opt/ai-tutor/deploy/backup/weekly-verify.sh" >> /var/log/ai-tutor-backup-verify.log 2>&1
